# supervisor.ini
[supervisord]
nodaemon = true
user = root

[program: server]
command = /bin/bash -c "exec /app/bin/server-linux-amd64"
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

[program: gateway]
command = /bin/bash -c "exec /app/bin/gateway-linux-amd64"
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

